# crontab.cron
@reboot /home/pi/lesbiancryptid/capture.sh


